# tendermonitor.service
[Unit]
Description=TenderMonitor - Мониторинг тендеров ЕИС
After=network.target nginx.service postgresql.service

[Service]
Type=simple
User=wanga
WorkingDirectory=/opt/tendermonitor
Environment="PATH=/opt/tendermonitor/venv/bin"
ExecStart=/opt/tendermonitor/venv/bin/python /opt/tendermonitor/main.py
Restart=always
RestartSec=30
TimeoutStartSec=3600
TimeoutStopSec=300
StandardOutput=journal
StandardError=journal

# Ограничение памяти: при превышении ~4 ГБ процесс будет перезапущен systemd
MemoryHigh=3G
MemoryLimit=4G

[Install]
WantedBy=multi-user.target

