# tendermonitor-status-migration.service
[Unit]
Description=TenderMonitor Daily Status Migration Service
After=network.target postgresql.service

[Service]
Type=oneshot
User=wanga
WorkingDirectory=/opt/tendermonitor
Environment="PATH=/opt/tendermonitor/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/tendermonitor/venv/bin/python3 -m database_work.daily_status_migration
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
